# model-switch.service
[Unit]
Description=Model Switch Panel
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/model-switch/model-switch.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
