# model-switch.service
[Unit]
Description=OpenClaw Model Switch Panel
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.openclaw/workspace/scripts/model-switch.py
Restart=on-failure
RestartSec=3
EnvironmentFile=%h/.openclaw/gateway.systemd.env

[Install]
WantedBy=default.target
